# plainnginx/supervisord/bootstrap.conf
[program:bootstrap-application]
command=/usr/local/bin/bootstrap.sh
priority=1
autorestart=false
startretries=0
startsecs=0
stdout_logfile= /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
