# docker/supervisord/bootstrap.conf
[program:bootstrap-application]
command=/usr/local/bin/bootstrap.sh
priority=1
autorestart=false
startretries=0
startsecs=0
